# init/wpdir.service
[Unit]
Description=WPDirectory - WordPress Plugin and Theme Directory Search
Documentation=https://github.com/wpdirectory/wpdir
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/wpdir
ExecStart=/var/www/wpdir/wpdir
Restart=always
RestartSec=10s
StandardOutput=append:/var/www/wpdir/wpdir.log
StandardError=append:/var/www/wpdir/wpdir.log

# Don't kill on timeout during startup (loading can take time)
TimeoutStartSec=infinity
TimeoutStopSec=30s

# Security settings
PrivateTmp=true
NoNewPrivileges=false

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target